# openclaw/rootfs/etc/s6-overlay/s6-rc.d/svc-gateway/run
#!/bin/sh
OPENCLAW_PORT="${OPENCLAW_PORT:-18789}"
exec openclaw gateway --port "${OPENCLAW_PORT}" --bind lan --allow-unconfigured
