# openclaw/rootfs/etc/s6-overlay/s6-rc.d/svc-playwright-mcp/run
#!/bin/sh
export DISPLAY=:99
# Wait for Xvfb to be ready (max 30s)
n=0; while ! xdpyinfo -display :99 >/dev/null 2>&1; do
  n=$((n+1)); [ "$n" -ge 60 ] && { echo "svc-playwright-mcp: Xvfb not ready after 30s" >&2; exit 1; }; sleep 0.5
done

PLAYWRIGHT_MCP_PORT="${PLAYWRIGHT_MCP_PORT:-8931}"

# Use Playwright's bundled Chromium if available, else fall back to system chromium
PW_CHROME=$(find "${PLAYWRIGHT_BROWSERS_PATH:-/opt/ms-playwright}" -name chrome -type f 2>/dev/null | head -1)
if [ -n "${PW_CHROME}" ]; then
    exec npx @playwright/mcp \
        --port "${PLAYWRIGHT_MCP_PORT}" \
        --host 0.0.0.0 \
        --executable-path "${PW_CHROME}" \
        --no-sandbox
else
    exec npx @playwright/mcp \
        --port "${PLAYWRIGHT_MCP_PORT}" \
        --host 0.0.0.0 \
        --executable-path /usr/bin/chromium-browser \
        --no-sandbox
fi
