# openclaw/rootfs/etc/s6-overlay/s6-rc.d/svc-novnc/run
#!/bin/sh
NOVNC_PORT="${NOVNC_PORT:-6080}"
exec /usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen "${NOVNC_PORT}"
